# systemd/calcurse-caldav.service
[Unit]
Description="sync calcurse with google calendar"

[Service]
ExecStart=calcurse-caldav --init=keep-remote
