# systemd/mailer.service
[Unit]
Description="push mail to remote"

[Service]
ExecStart=/home/ark/.local/bin/mailer push
