# systemd/popup-killer.service
[Service]
Type=simple
ExecStart=/home/retroflexivity/bin/popup-killer

[Install]
WantedBy=default.target
