# systemd/power-control.service
[Service]
ExecStart=/home/ark/.local/bin/power-control start
Restart=always
RestartSec=3s

[Install]
WantedBy=default.target
