# systemd/power-mode-switcher.service
[Service]
Environment=STARTUP_WAIT=30s
ExecStart=/home/ark/bin/power-mode-switcher

[Install]
WantedBy=default.target
