# systemd/rclone-gdrive.service
[Unit]
Description="rclone bisync google drive"

[Service]
ExecStart=/home/ark/.local/bin/rclone-gdrive
