# systemd/rclone-gdrive.service
[Unit]
Description="rclone bisync google drive"

[Service]
ExecStart=/home/retroflexivity/bin/rclone-gdrive
